# shadowsock/resource/shadowsocks.conf
[program:shadowsocks]
command=ssserver -c /etc/shadowsocks.json
autorestart=true
autostart=true
startsecs=3
user=nobody
